# launchd/com.dcr.fub-refresh.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  FUB daily refresh. Install:
    cp launchd/com.dcr.fub-refresh.plist ~/Library/LaunchAgents/
    launchctl load ~/Library/LaunchAgents/com.dcr.fub-refresh.plist
  Schedules are offset 30 min from the GHL job (:00) so both jobs don't open
  connections to the shared VM at the same instant.
-->
<plist version="1.0">
<dict>
  <key>Label</key><string>com.dcr.fub-refresh</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/smokestack/Projects/FUB_API/.venv/bin/python</string>
    <string>/Users/smokestack/Projects/FUB_API/scripts/refresh_daily.py</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/smokestack/Projects/FUB_API</string>

  <!-- 9:30 / 13:30 / 16:30 / 21:30 local — offset 30 min from GHL's :00 job -->
  <key>StartCalendarInterval</key>
  <array>
    <dict><key>Hour</key><integer>9</integer><key>Minute</key><integer>30</integer></dict>
    <dict><key>Hour</key><integer>13</integer><key>Minute</key><integer>30</integer></dict>
    <dict><key>Hour</key><integer>16</integer><key>Minute</key><integer>30</integer></dict>
    <dict><key>Hour</key><integer>21</integer><key>Minute</key><integer>30</integer></dict>
  </array>

  <!-- Don't run at load; missed fires coalesce into one catch-up run on wake -->
  <key>RunAtLoad</key><false/>
  <key>StartCalendarIntervalCausesWake</key><false/>

  <key>StandardOutPath</key>
  <string>/Users/smokestack/Projects/FUB_API/logs/refresh-out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/smokestack/Projects/FUB_API/logs/refresh-err.log</string>
</dict>
</plist>
